# gc2conflict/conf/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:gc2conflict]
command=/bin/bash -c "nodejs ~/gc2conflict/app/app.js -D FOREGROUND"

[program:html2pdf]
command=/bin/bash -c "nodejs ~/html2pdf.it -D FOREGROUND"
