# logstash/apache/conf/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:logstash]
command=/bin/bash -c "/usr/share/logstash/bin/logstash -f /root/logstash-apache.conf"
priority=1
redirect_stderr=true

[program:nodejs]
command=/bin/bash -c "cd /root/server && npm start"
priority=2
redirect_stderr=true